# .travis/config/supervisord.conf
[supervisord]
nodaemon=true

[program:startup]
priority=1
command=service salt-minion start
stdout_logfile=/var/log/salt/minion
stderr_logfile=/var/log/salt/minion
